# scripts/backup.service
[Unit]
Description=Script to backup game files
Wants=timers.target

[Service]
Type=simple
ExecStart=/home/mc/bedrockServer/scripts/backup.sh
